# dnscrypt/root/etc/services.d/dnscrypt-proxy/run
#!/command/execlineb -P

/usr/sbin/dnscrypt-proxy --config /etc/dnscrypt-proxy/dnscrypt-proxy.toml
